# pybluetooth.conf
description "Bluetooth Server"

start on started /etc/init.d/bluetooth

exec sh /home/pi/bot.sh
